# cecfix.service
[Unit]
Description=CEC Fix
After=multi-user.target

[Service]
Type=idle
ExecStart={{DIR}}/build/cec-fix
Restart=always
RestartSec=5
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target